# crontab
# Set PATH for cron environment
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Run the Python script every minute
* * * * * cd /app && /usr/local/bin/python3 /app/src/main.py >> /proc/1/fd/1 2>/proc/1/fd/2
